# ops/linux/services/moza_udp_client.service.in
[#] Template unit file for RC Simulator (MOZA UDP Client)
[#] This file is not installed directly. Use:
[#]   ops/linux/install_service.sh --user <username>
[#]
[#] Placeholders:
[#]   @USER@        system user to run the service as
[#]   @DISPLAY@     X/Wayland display value (default :0)
[#]   @XAUTHORITY@  XAUTHORITY path (default %h/.Xauthority)
[#]   @PYTHON@      python executable path (default /usr/bin/python3)

[Unit]
Description=RC Simulator (MOZA UDP Client)
After=network.target
Wants=network.target

[Service]
Type=simple
User=@USER@
WorkingDirectory=%h
Environment=DISPLAY=@DISPLAY@
Environment=XAUTHORITY=@XAUTHORITY@
Environment=PYTHONUNBUFFERED=1
# Prefer the console script when @PYTHON@ points to a venv, but keep module fallback.
ExecStart=/bin/bash -lc 'PYDIR="$(dirname "@PYTHON@")"; if [ -x "${PYDIR}/rc-simulator-headless" ]; then exec "${PYDIR}/rc-simulator-headless"; else exec "@PYTHON@" -m rc_simulator.__main_headless__; fi'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target

